# systemd/admin-web.service
[Unit]
Description=VK Auto Bot Admin Web (FastAPI/uvicorn)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=vkbot
Group=vkbot
WorkingDirectory=/opt/vk_auto_bot
ExecStart=/opt/vk_auto_bot/.venv/bin/uvicorn app.admin.main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=127.0.0.1 --no-server-header
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
