# systemd/vk-bot.service
[Unit]
Description=VK Auto Bot (Long Poll + worker)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=vkbot
Group=vkbot
WorkingDirectory=/opt/vk_auto_bot
ExecStart=/opt/vk_auto_bot/.venv/bin/python -m app.bot.main
# on-failure (не always): без VK-токена бот штатно завершается (exit 0) и не уходит
# в бесконечный цикл рестартов; реальный сбой (exit != 0) — перезапустит.
# Токен задаётся в админке → после этого: systemctl restart vk-bot
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
